# systemd/survival-recovery.service
[Unit]
Description=Remote Survival Kit - Recovery Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/remote-survival/recovery_monitor.sh